# cron.d/low-battery-handler
* * * * * root /usr/local/bin/low-battery-handler
